# webapp/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log;
pidfile=/var/run/supervisord.pid;
loglevel=debug

[program:kms]
command=/bin/bash /kms.sh
redirect_stderr=true

[program:openvidu-server]
command=/bin/bash /openvidu.sh
redirect_stderr=true

[program:app]
command=/bin/bash -c "http-server -S -C /webapp/cert.pem -K /webapp/key.pem /webapp/"
redirect_stderr=true